# .claude/worktrees/agent-a96de9bd/Pryan-Fire/hughs-forge/scripts/position-monitor.service
[Unit]
Description=Position Monitor — DLMM position feed to Discord
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/data/repos/The-Nexus/Pryan-Fire
ExecStart=/data/repos/The-Nexus/Pryan-Fire/hughs-forge/services/trade-orchestrator/venv/bin/python hughs-forge/scripts/position_monitor.py
StandardOutput=journal
StandardError=journal
