# Pryan-Fire/hughs-forge/services/trade-orchestrator/patryn-trader.service
[Unit]
Description=Patryn Trade Orchestrator
After=network.target

[Service]
Type=simple
WorkingDirectory=/data/openclaw/workspace/Pryan-Fire
ExecStart=/data/openclaw/workspace/Pryan-Fire/hughs-forge/services/trade-orchestrator/start-patryn-trader.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
