# Pryan-Fire/hughs-forge/services/trade-orchestrator/patryn-trader.service
[Unit]
Description=Patryn Trade Orchestrator
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/data/repos/Pryan-Fire
ExecStart=/data/openclaw/workspace/Pryan-Fire/hughs-forge/services/trade-orchestrator/start-patryn-trader.sh
Restart=always
RestartSec=10
StandardOutput=append:/var/log/patryn-trader.log
StandardError=append:/var/log/patryn-trader.log

[Install]
WantedBy=default.target
